# conf/h/.config/autostart/setup-complete.desktop
[Desktop Entry]
Name=FinishUserConfiguration
Exec=~/setup-complete.sh
Terminal=true
Type=Application
Encoding=UTF-8
Categories=System;Application;