# conf/h/.config/autostart/skype.desktop
[Desktop Entry]
Type=Application
Name=skype
Exec=/usr/bin/skype
