# conf/h/.config/autostart/volumeicon.desktop
[Desktop Entry]
Name=volume Widget
Exec=volumeicon
Terminal=false
Type=Application
Encoding=UTF-8
Categories=System;Application;