# conf/h/.config/autostart/xscreensaver.desktop
[Desktop Entry]
Name=xscreensaver
Exec=xscreensaver -no-splash
Terminal=false
Type=Application
Encoding=UTF-8
Categories=System;Application;
